# airootfs/etc/systemd/system/zeronet.service.d/stop-quicker.conf
[Service]
## Zeronet takes extremely long to stop usually so
## this makes systemd send a SIGKILL signal to stop it
## which kills it instantly.
KillSignal=SIGKILL
