# airootfs/etc/systemd/system/zeronet.service.d/wait-for-tor.conf
[Unit]
# Wait til Tor has started before running.
After=tor.service
Requires=tor.service
After=secure-time-sync.service
Requires=secure-time-sync.service
