# airootfs/etc/profile.d/umask.sh
# Set a more restrictive umask.
# This makes newly created files only readable by the owner.
umask 0077
