# airootfs/etc/systemd/system/NetworkManager.service.d/fail-closed.conf
[Unit]
## Fail Closed Mechanism.
## If iptables isn't running then NetworkManager won't start.
After=iptables.service
Requires=iptables.service
